# deploy/oracle/dsa.service
[Unit]
Description=Daily Stock Analysis
After=network.target docker.service
Wants=docker.service

[Service]
Type=simple
WorkingDirectory=/opt/daily_stock_analysis
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/daily_stock_analysis/.venv/bin/python -m uvicorn api.app:app --host 0.0.0.0 --port 8080
Restart=always
RestartSec=5
User=root

[Install]
WantedBy=multi-user.target
